# config/systemd/user/distrobox-upgrade.service
[Unit]
Description=Distrobox containers automatic updates

[Service]
Type=simple
ExecStart=distrobox-upgrade --all
StandardOutput=null
