# config/systemd/user/distrobox-update.timer
[Unit]
Description=Distrobox containers automatic updates trigger

[Timer]
OnUnitInactiveSec=1d
OnBootSec=10min

[Install]
WantedBy=timers.target
