# systemd/user/distrobox-update.service
[Unit]
Description=Distrobox containers automatic updates

[Service]
Type=oneshot
ExecStart=distrobox-upgrade --all
StandardOutput=null
